# daemon/qubes-db.socket
[Socket]
ListenStream=/run/qubes/qubesdb.sock

[Install]
WantedBy=sockets.target
